# scripts/com.daedal.kidsnews-autofix.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.daedal.kidsnews-autofix</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <!-- {{HOME}} placeholder is replaced by install-autofix-daemon.sh -->
        <string>{{HOME}}/myprojects/news-v2/scripts/autofix-daemon.sh</string>
    </array>

    <!-- Fire ONCE/day at 04:00 local time (Eastern, since the owner's
         Mac is set to ET). This window picks the lowest possible
         interference: owner is asleep, so back-to-back `claude -p`
         spawns can't lock their IDE. The daemon drains ALL queued
         items in one pass, not one-per-tick.

         Calendar-based scheduling beats StartInterval because modern
         macOS suspends launchd jobs during sleep, so an interval
         "every 24h" drifts unpredictably. With an explicit HH:MM, a
         missed tick (lid closed at 04:00, opened at 09:00) fires
         catch-up at next wake instead of silently disappearing.

         Routine quality fixes run in CI (autofix_apply.py); this
         daemon ONLY drains rows the admin explicitly clicked
         "🤖 Fix with Claude" on in a digest email — i.e.
         status='fix-requested'. Most days the queue is empty and
         the daemon exits within seconds. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>4</integer>
        <key>Minute</key><integer>0</integer>
    </dict>

    <!-- Run a single tick immediately at load time so the user can
         verify it's working without waiting for the next 4-hour
         window. -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Don't auto-restart on success exit. -->
    <key>KeepAlive</key>
    <false/>

    <!-- launchd consolidates these to ~/Library/Logs/launchd_*. The
         daemon writes its own structured log under
         ~/Library/Logs/kidsnews-autofix/ — these two files are just
         for catastrophic launchd-level failures. -->
    <key>StandardOutPath</key>
    <string>{{HOME}}/Library/Logs/kidsnews-autofix-out.log</string>
    <key>StandardErrorPath</key>
    <string>{{HOME}}/Library/Logs/kidsnews-autofix-err.log</string>

    <!-- Inherit user's $PATH so node, python, claude, etc. resolve. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>{{HOME}}/.npm-global/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
