# scripts/com.daedal.kidsnews-autofix.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.daedal.kidsnews-autofix</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <!-- {{HOME}} placeholder is replaced by install-autofix-daemon.sh -->
        <string>{{HOME}}/myprojects/news-v2/scripts/autofix-daemon.sh</string>
    </array>

    <!-- Run every 8 hours (3 ticks per day). Owner picked this
         cadence — autofix isn't time-sensitive (pipeline runs once
         a day, fixes can settle within hours), and 8h keeps the
         daemon nearly silent in logs while still draining the queue
         within 24h of a problem appearing. The daemon exits within
         seconds when the queue is empty. -->
    <key>StartInterval</key>
    <integer>28800</integer>

    <!-- Run a single tick immediately at load time so the user can
         verify it's working without waiting for the next half-hour. -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Don't auto-restart on success exit. -->
    <key>KeepAlive</key>
    <false/>

    <!-- launchd consolidates these to ~/Library/Logs/launchd_*. The
         daemon writes its own structured log under
         ~/Library/Logs/kidsnews-autofix/ — these two files are just
         for catastrophic launchd-level failures. -->
    <key>StandardOutPath</key>
    <string>{{HOME}}/Library/Logs/kidsnews-autofix-out.log</string>
    <key>StandardErrorPath</key>
    <string>{{HOME}}/Library/Logs/kidsnews-autofix-err.log</string>

    <!-- Inherit user's $PATH so node, python, claude, etc. resolve. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>{{HOME}}/.npm-global/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
