# scripts/com.daedal.kidsnews-autofix.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.daedal.kidsnews-autofix</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <!-- {{HOME}} placeholder is replaced by install-autofix-daemon.sh -->
        <string>{{HOME}}/myprojects/news-v2/scripts/autofix-daemon.sh</string>
    </array>

    <!-- Fire 4×/day at fixed local times (06:00, 10:00, 14:00, 18:00).
         Calendar-based scheduling beats StartInterval because modern
         macOS suspends launchd jobs during sleep, so a "every 4h"
         interval drifts unpredictably overnight. With explicit
         HH:MM, missed ticks (lid closed at the fire time) reschedule
         to the next listed time, never silently disappear.

         Routine quality fixes run in CI now (autofix_apply.py); this
         daemon ONLY drains rows the admin explicitly clicked
         "🤖 Fix with Claude" on in a digest email — i.e.
         status='fix-requested'. Most days the queue is empty and
         the daemon exits within seconds. -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key><integer>6</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>10</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>14</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>18</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
    </array>

    <!-- Run a single tick immediately at load time so the user can
         verify it's working without waiting for the next 4-hour
         window. -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Don't auto-restart on success exit. -->
    <key>KeepAlive</key>
    <false/>

    <!-- launchd consolidates these to ~/Library/Logs/launchd_*. The
         daemon writes its own structured log under
         ~/Library/Logs/kidsnews-autofix/ — these two files are just
         for catastrophic launchd-level failures. -->
    <key>StandardOutPath</key>
    <string>{{HOME}}/Library/Logs/kidsnews-autofix-out.log</string>
    <key>StandardErrorPath</key>
    <string>{{HOME}}/Library/Logs/kidsnews-autofix-err.log</string>

    <!-- Inherit user's $PATH so node, python, claude, etc. resolve. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>{{HOME}}/.npm-global/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
